# services/systemd_files/power_prediction.service
[Unit]
Description="power prediction service"
Wants=power_prediction.timer

[Service]
Type=oneshot
EnvironmentFile=/etc/ems/EMS_systemd_config.txt
ExecStart=/usr/bin/bash -c 'source "$BASHRCPATH"; cd "$EMSFOLDER"; conda activate milp; python -m services.power_prediction'

[Install]
WantedBy=multi-user.target                         