# services/systemd_files/power_prediction.service
[Unit]
Description="remove the old user data from the currently in use folder"
Wants=power_prediction.timer

[Service]
Type=oneshot
Environment=BASHRCPATH= EMSFOLDER=
ExecStart=/usr/bin/bash -c 'source "$BASHRCPATH"; cd "$EMSFOLDER"; conda activate milp; python -m services.learning_cleaner'

[Install]
WantedBy=multi-user.target                         