# services/systemd_files/power_prediction.timer
[Unit]
Description="timer for the power prediction service"
Requires=power_prediction.service
[Timer]
Unit=power_prediction.service
OnCalendar=*-*-* *:14/15:00
AccuracySec=20s

[Install]
WantedBy=timers.target